# bin/Release/minimessanger/mm.service
[Unit]
Description=Minimessanger server
After=multi-user.target
Requires=mysql.service

[Service]
WorkingDirectory=/var/www/minimessanger
ExecStart=sudo /var/www/minimessanger/mm
Restart=always

[Install]
WantedBy=multi-user.target

